# user-devenv/devenv/dockerfiles/phpserver/etc/service/php5-fpm/run
#!/bin/bash
exec 2>&1

. /etc/container_environment.sh
/usr/sbin/php5-fpm -F -c "/etc/php5/fpm/php.ini.d/$V3_ENV/$V3_ROLE-php.ini"
